# monitoring.service
[Unit]
Description=Simple Process Monitoring
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring.sh
User=root
